# crontab
* * * * * cd /usr/src/app && ./cron-run-script.sh

# Don't forget to end your cron file with an empty line
